# services/mitm-lab-suricata.service
[Unit]
Description=MITM lab live Suricata sensor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mitm-lab-suricata-live.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
