# services/mitm-lab-zeek.service
[Unit]
Description=MITM lab live Zeek sensor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mitm-lab-zeek-live.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
