# services/mitm-lab-detector.service
[Unit]
Description=MITM lab victim detector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mitm_lab_detector.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
